# project/final_project/parmco-ble.service
# ==============================================================================
# PARMCO BLE Server — systemd Service Unit
# ==============================================================================
#
# Project:  PARMCO (Phone APP RP4 Motor Control)
# Course:   ECSE 4235 — Embedded Systems, Spring 2026
# Authors:  Justin Moreno and Jalen Edusei
#
# AI Attribution:
#   This service file was AI-Generated by Claude Opus 4.6 (Anthropic).
#
# Purpose:
#   Starts the PARMCO BLE GATT server automatically on boot so the RP4
#   is ready to accept phone connections without any manual intervention.
#   The project spec requires: "The RP4 must boot and run the code."
#
# Install:
#   sudo cp parmco-ble.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable parmco-ble.service
#   sudo systemctl start parmco-ble.service
#
# Management:
#   sudo systemctl status parmco-ble    — check if running
#   sudo systemctl restart parmco-ble   — restart the server
#   sudo systemctl stop parmco-ble      — stop the server
#   sudo journalctl -u parmco-ble -f    — view live logs
#
# ==============================================================================

[Unit]
# Human-readable description shown in systemctl status output
Description=PARMCO BLE Motor Control Server (ECSE 4235)

# Wait for the BlueZ daemon itself (not just bluetooth.target) so our D-Bus
# calls to org.bluez succeed. CP3 spec F6: Requires= not Wants=, and bind to
# bluetooth.service so we restart when BlueZ restarts.
After=bluetooth.service
Requires=bluetooth.service

[Service]
# Type=simple means systemd considers the service started as soon as
# the process is forked. This is correct for our GLib main loop program.
Type=simple

# Run as root — required for D-Bus system bus access and BlueZ interaction.
# A more restrictive approach would use a dedicated user with a custom
# D-Bus policy file, but root is simpler for a course project.
User=root

# Full path to the installed binary
ExecStart=/usr/local/bin/ble_server

# On failure, wait 3 seconds and restart automatically.
# This handles transient BlueZ startup race conditions.
Restart=on-failure
RestartSec=3

# Allow up to 5 restart attempts within 30 seconds before giving up
StartLimitIntervalSec=30
StartLimitBurst=5

# Send stdout/stderr to the systemd journal for debugging.
# View with: sudo journalctl -u parmco -f
StandardOutput=journal
StandardError=journal

# Give the process up to 5 seconds to shut down cleanly on stop/restart.
# Our SIGTERM handler calls g_main_loop_quit() and cleans up D-Bus objects.
TimeoutStopSec=5

[Install]
# Start this service when the system reaches multi-user mode (normal boot).
WantedBy=multi-user.target

